# infra/supervisor/supervisord.conf
[supervisord]
user=root
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0
pidfile=/tmp/supervisord.pid

[program:intelligence]
command=uvicorn main:app --host 127.0.0.1 --port 8000
directory=/app/services/ai-runtime
autostart=true
autorestart=true
startsecs=5
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=PYTHONUNBUFFERED="1",INTELLIGENCE_HOST="127.0.0.1",AI_PORT="8000"

[program:platform]
command=node dist/index.js
directory=/app/services/gateway
autostart=true
autorestart=true
startsecs=5
priority=20
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=NODE_ENV="production",API_PORT="%(ENV_PORT)s"
